# graphical-session/.config/systemd/user/graphical-session@.target
[Unit]
Description=Graphical session at '%i'
