# graphical-session/.config/systemd/user/xsession@.target
[Unit]
Description=Xsession on DISPLAY=:%i
