# foot/.config/systemd/user/foot-server@.socket
[Unit]
Requisite=wayland-session@%i.target

[Socket]
ListenStream=%t/foot-%i.sock

[Install]
WantedBy=wayland-session@%i.target
