# xsession/.config/systemd/user/xinit-remaps@.service
[Unit]
Description=Oneshot xodmap/setxkbmap config at graphical session '%i' 
Documentation=man:setxkbmap(1) man:xmodmap(1)
#PartOf=xorg@%i.target
After=graphical-session@%i.target
Requisite=graphical-session@%i.target

[Service]
Type=oneshot
ExecStart=sh -c 'session -s "%i" run "xinit-remaps"'
; Restart=on-failure

[Install]
WantedBy=graphical-session@%i.target
