# xsession/.config/systemd/user/xrandr-tearfree@.service
[Unit]
Description=Oneshot xrandr TearFree to 'on' at graphical session '%i' 
Documentation=man:xrandr(1)
#PartOf=xorg@%i.target
After=graphical-session@%i.target
Requisite=graphical-session@%i.target

[Service]
Type=oneshot
ExecStart=sh -c 'session -s "%i" run "xrandr --output HDMI-A-1 --set TearFree on"'
; Restart=on-failure

[Install]
WantedBy=graphical-session@%i.target
