# xsession/.config/systemd/user/xsession@.target
[Unit]
Description=X graphical session on DISPLAY=:%i
