# wayland-session/.config/systemd/user/wayland-session@.target
[Unit]
Description=Wayland graphical session on WAYLAND_DISPLAY=%i
